# src/jdk.jpackage/linux/classes/jdk/jpackage/internal/resources/unit-template.service
[Unit]
Description=APPLICATION_DESCRIPTION

[Service]
ExecStart=APPLICATION_LAUNCHER
Restart=on-failure

[Install]
WantedBy=multi-user.target
